# services/plex.service.conf
[program:plex]
user=docker
command=/usr/local/bin/start-plex.sh
priority=0
stopwaitsecs=5
stdout_events_enabled = true
stderr_events_enabled = true
